# src/Roberta.Hub/Roberta.Hub/robhub.service
[Unit]
Description=robcli service
After=network.target

[Service]
Environment=DOTNET_ROOT=/home/pi/.dotnet
Environment=PATH=$PATH:/home/pi/.dotnet'
WorkingDirectory=/home/pi/roberta/hub
ExecStart=/home/pi/roberta/hub/robchub
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
